# ai.tailscale.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.tailscale.serve</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Tailscale.app/Contents/MacOS/Tailscale</string>
        <string>serve</string>
        <string>--bg</string>
        <string>4096</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/mihai</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mihai</string>
    </dict>
    
    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep running if it crashes -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- Restart after failure with 10 second delay -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/tmp/tailscale-serve.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tailscale-serve.error.log</string>
</dict>
</plist>